# build/app/supervisor/app.conf
[supervisord]
nodaemon=true

[program:app]
command = python -m application.coinbase.feed
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
